# scripts/com.fantasy-baseball.fantrax-pickup-audit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.fantasy-baseball.fantrax-pickup-audit</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/emet_macbook_air/Documents/Fantasy Baseball Analytics/scripts/run_fantrax_pickup_audit.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/emet_macbook_air/Documents/Fantasy Baseball Analytics</string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/emet_macbook_air/Documents/Fantasy Baseball Analytics/outputs/fantrax_export/fantrax_pickup_audit_launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/emet_macbook_air/Documents/Fantasy Baseball Analytics/outputs/fantrax_export/fantrax_pickup_audit_launchd.err</string>
</dict>
</plist>
